# src/hardware/systemd/robot-console.service
[Unit]
Description=Robot Control Web Console (RPI-01)
Documentation=file:///home/pi/Desktop/web_console.py
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
SupplementaryGroups=dialout audio video input
WorkingDirectory=/home/pi/Desktop
Environment=PYTHONUNBUFFERED=1
# 用 yolo_env venv 的 python (包含 ultralytics + cv2 + torch + fastapi + uvicorn + pyserial)
ExecStart=/home/pi/yolo_env/bin/python3 /home/pi/Desktop/web_console.py
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
